# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:hive]
command=/opt/hive/bin/hive --service hiveserver2
autostart=true
autorestart=false
numprocs=1
